# cron_service/system/crontab
# START CRON JOB
* * * * * /usr/bin/env -i $(cat /proc/1/environ | tr '\0' ' ') /usr/local/bin/python /app/src/main.py > /proc/1/fd/1 2>/proc/1/fd/2
# END CRON JOB
